# desktop/services/conf.d/xfce4.sv.conf
[program:xfce4]
environment=HOME="/home/alpine",DISPLAY=":1",USER="alpine"
command=/usr/bin/startxfce4
user=alpine
autorestart=true
priority=300
